# eora/deploy/eora.service
# Пример: /etc/systemd/system/eora.service
# Пользователь www-data — убедитесь, что у него есть чтение /srv/eora и запись в media при необходимости.

[Unit]
Description=EORA Django (Gunicorn)
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/srv/eora
EnvironmentFile=/etc/eora.env
ExecStart=/srv/eora/.venv/bin/gunicorn eora.wsgi:application \
  --bind 127.0.0.1:8000 \
  --workers 3 \
  --timeout 120
Restart=always

[Install]
WantedBy=multi-user.target
